# recipes-ynov/sd-bus-exemple/sd-bus-example/sd-server.service
[Unit]
Description=Sd Server service
[Service]
Type=simple
ExecStart=/usr/sbin/sd-server
[Install]
WantedBy=multi-user.target